# tools/frpc/etc/supervisord.conf
; Debian 3.17 supervisord
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/run/supervisord.pid

[inet_http_server]
port = 0.0.0.0:9001
username = Kylin
password = Kylin

[program:sshd]
command=/usr/sbin/sshd -D
user=root
autostart=true
autorestart=true

[program:crond]
command=/usr/sbin/crond -f
user=root
autostart=true
autorestart=true

[program:frpc]
command=/usr/local/bin/frpc -c /etc/frp/frpc.ini
user=root
autostart=true
autorestart=true

[program:v2ray]
command=/usr/local/bin/kylin-v2ray
user=root
autostart=true
autorestart=true

[program:nginx]
command=nginx -c /data/wwwconf/nginx.conf -g 'daemon off;'
user=root
autostart=true
autorestart=true
