# tools/satis/etc/supervisord.conf
; satis supervisord
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/run/supervisord.pid

[inet_http_server]
port = 0.0.0.0:9001
username = Kylin
password = Kylin

[program:sshd]
command=/usr/sbin/sshd -D
user=root
autostart=true
autorestart=true

[program:crond]
command=/usr/sbin/crond -f
user=root
autostart=true
autorestart=true

[program:satis]
command=/satis/bin/docker-entrypoint.sh --ansi --ansi build /build/satis.json /build/output
user=root
autostart=false
autorestart=false
